# package/network/services/atftp/files/atftpd.init
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org

START=98
STOP=15
USE_PROCD=1
ATFTPD_BIN="/usr/sbin/atftpd --user root --daemon /tmp"

start_service()
{
	[ -f /tmp/bob_cm.cfg ] || {
		cp -f /etc/bob_cm.cfg /tmp/bob_cm.cfg
		chmod 777 /tmp/bob_cm.cfg
	}
	procd_open_instance
	procd_set_param respawn
	procd_set_param command $ATFTPD_BIN
	procd_close_instance
}